# ai-roleplay-site.service
# AI Roleplay Site systemd service
# 作用：以 systemd 管理最简 AI 角色扮演网站，便于开机自启与统一启停。

[Unit]
Description=AI Roleplay Site
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/openclaw/workspace/project/app/louge-roleplay
Environment=NODE_ENV=production
ExecStart=/root/.nvm/versions/node/v24.15.0/bin/node /opt/openclaw/workspace/project/app/louge-roleplay/src/server.js
Restart=always
RestartSec=5
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
